# deploy/qortium-network-collect.service
[Unit]
Description=Qortium Network topology collector (capture one snapshot)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/qortium-network
EnvironmentFile=%h/qortium-network/deploy/qortium-network.env
# Userland Node (see deploy/README.md). Adjust if installed elsewhere.
ExecStart=%h/.local/node/bin/node scripts/auto-collect.mjs
Nice=10
