# deploy/qortium-network-publish.service
[Unit]
Description=Qortium Network topology publisher (post best recent record)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/qortium-network
EnvironmentFile=%h/qortium-network/deploy/qortium-network.env
# Userland Node (see deploy/README.md). Adjust if installed elsewhere.
ExecStart=%h/.local/node/bin/node scripts/auto-publish.mjs
Nice=10
